# services.d/gpg/finish
#!/bin/sh
echo "[END]  gpg with error: $1" >&2

sleep 5
exit $1
