# services.d/pcscd/finish
#!/bin/sh
echo "[END]  pcscd with error: $1" >&2

sleep 5
exit $1
