# services.d/pcscd/run
#!/bin/sh
s6-svwait -U /var/run/s6/services/init

echo "[START] pcscd"
env

exec /usr/sbin/pcscd --info --foreground --color --disable-polkit
